# systemd/mnt-shared.mount
[Mount]
What=UUID=37584AA3704AC8C1
Where=/mnt/shared
Type=ntfs-3g
Options=permissions

[Install]
WantedBy=local-fs.target
